# xRootD-TPC-GridFTP/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
childlogdir = /var/log/supervisor

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket
loglevel=debug

[program:initconfig]
command=/usr/local/sbin/fix_certs.sh
autorestart=false
startretries=0
startsecs=0
redirect_stderr=false
redirect_stderr=false

[program:gridftp]
command=/usr/sbin/globus-gridftp-server -c /etc/gridftp.conf -C /etc/gridftp.d -pidfile /var/run/globus-gridftp-server.pid -aa
autorestart=true

[program:crond]
command=/usr/sbin/crond -n
autorestart=true
