# maddash/10-httpd.conf
[program:httpd]
command=/usr/sbin/httpd
autorestart=true